# build/files/supervisord.conf
[unix_http_server]
file=/run/supervisord.sock   ; (the path to the socket file)
user=root

[supervisord]
logfile=/dev/stdout
logfile_maxbytes=0
loglevel=debug                ; (log level;default info; others: debug,warn,trace)

[inet_http_server]         ; inet (TCP) server disabled by default
port=127.0.0.1:9001        ; ip_address:port specifier, *:port for all iface

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///run/supervisord.sock 

# [program:monitor]
# directory=/usr/src/monitor
# command=npm start --production "/data"
# autostart=true
# redirect_stdout=true
# redirect_stderr=true

[program:miner]
command=/opt/miner/bin/miner foreground
autostart=true
redirect_stdout=true
redirect_stderr=true
